# systemd/sandd.socket
[Unit]
Description=Countdown timer daemon socket

[Socket]
ListenStream=%t/sandd.sock
SocketMode=0666

[Install]
WantedBy=sockets.target
