# systemd/sandd.service
[Unit]
Description=Countdown timer daemon
Requires=sandd.socket

[Service]
ExecStart=/usr/bin/sandd

[Install]
WantedBy=default.target
