# resources/systemd/sand.socket
[Unit]
Description=Countdown timer daemon socket

[Socket]
ListenStream=%t/sand.sock
SocketMode=0666

[Install]
WantedBy=sockets.target
